# crontab
0 10,16 * * * python /usr/src/ssl_expiry_check/main.py >> /var/log/cron/cron.log 2>&1
